# cron.d/hashes
#!/bin/sh

/usr/local/sbin/vtpipeline hashes
/usr/local/sbin/vtpipeline hashes --date-offset 90
