# cron.d/scans
#!/bin/sh

/usr/local/sbin/vtpipeline scans
/usr/local/sbin/vtpipeline scans --date-offset 90
